# .config/systemd/arch-update.service
[Unit]
Description=Automated Arch Linux System Update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/%u/.config/i3/arch-update.sh
StandardOutput=journal
StandardError=journal

# Security settings
PrivateTmp=yes
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=/var/log

[Install]
WantedBy=multi-user.target